# ops/routine/com.givyx.ops-routine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Ops answer-pickup routine (OpsPA §8). Runs hourly while the Mac is awake. Install:
       cp ops/routine/com.givyx.ops-routine.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.givyx.ops-routine.plist
     Uninstall: launchctl unload ~/Library/LaunchAgents/com.givyx.ops-routine.plist -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.givyx.ops-routine</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/stan/Code/givyx/PersonalAssistant/ops/routine/pickup-answers.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/stan/.givyx/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/stan/.givyx/launchd.err.log</string>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
